# ecolearn_backend/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/usr/sbin/mysqld
autostart=true
autorestart=true
priority=5

[program:spring]
command=java -jar /app/app.jar
autostart=true
autorestart=true
priority=10
